# systemd-user/.config/systemd/user/aria2.service
[Unit]
Description=Aria2 download daemon

[Service]
ExecStart=/usr/bin/aria2c

[Install]
WantedBy=multi-user.target
