# systemd-user/.config/systemd/user/mail_sync.timer
[Unit]
Description=Main sync timer

[Timer]
Persistent=true
OnCalendar=*:0/30

[Install]
WantedBy=timers.target