# systemd-user/.config/systemd/user/mbsync.timer
[Unit]
Description=Mail sync timer

[Timer]
OnCalendar=*-*-* 00/1:00:00
Persistent=true

[Install]
WantedBy=timers.target
