# systemd-user/.config/systemd/user/synergyc@.service
[Unit]
Description=Synergy Client Daemon for %i
After=network.target

[Service]
ExecStart=/usr/bin/synergyc --no-daemon %i
Restart=on-failure
